# ansible/mk_copy_cat/files/copy_cat.timer
[Unit]
Description=copy_cat timer

[Timer]
OnCalendar=*-*-* 07,15,23:00:00

[Install]
WantedBy=timers.target
