# ansible/mk_copy_bd/files/copy_bd.timer
[Unit]
Description=copy_db timer

[Timer]
OnCalendar=*-*-* 07,15,23:00:00

[Install]
WantedBy=timers.target
